# bulbs.service
[Unit]
Description=bulbs.py Server für myStrom-Birnen
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/etc/network/bulbs.py
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=default.target
